# system/logfalcon-firstboot.service
# systemd unit for logfalcon boot-time config application
# Reads /boot/firmware/logfalcon-config.txt and updates hostapd + app config.
# Install to: /etc/systemd/system/logfalcon-firstboot.service
# Enable: sudo systemctl enable logfalcon-firstboot.service

[Unit]
Description=LogFalcon Boot Config
Documentation=https://github.com/proeugene/logfalcon
After=local-fs.target
Before=hostapd.service logfalcon-web.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/logfalcon/firstboot.sh

[Install]
WantedBy=multi-user.target
