# system/logfalcon@.service
# systemd template unit for LogFalcon sync service
# Install to: /etc/systemd/system/logfalcon@.service
# Triggered by udev when FC is plugged in (e.g. logfalcon@ttyACM0.service)

[Unit]
Description=LogFalcon (%I)
Documentation=https://github.com/proeugene/logfalcon
# Stop if the device disappears mid-sync
BindsTo=dev-%i.device
After=dev-%i.device network.target
# Only one sync at a time
Conflicts=logfalcon@*.service

[Service]
Type=oneshot
# Wait for USB enumeration to stabilise
ExecStartPre=/bin/sleep 3

User=bbsyncer
Group=dialout

# Install location
WorkingDirectory=/opt/logfalcon
ExecStart=/opt/logfalcon/venv/bin/python -m logfalcon --port /dev/%I

StandardOutput=journal
StandardError=journal
SyslogIdentifier=logfalcon

# Resource limits
TimeoutStartSec=600
TimeoutStopSec=10

# Do not restart on failure — the pilot will re-plug the FC
Restart=no

[Install]
WantedBy=multi-user.target
